# systemd/legoarm.service
[Unit]
Description=LEGO Arm HTTP API (stdlib)
After=network-online.target
Wants=network-online.target

[Service]
User=lego
Group=lego
WorkingDirectory=/path/to/lego-arm
ExecStart=/usr/bin/python3 /path/to/lego_arm_master.py
Environment=PORT=5001
Environment=API_KEY=<YOUR_API_KEY>
Environment=USE_FAKE_MOTORS=1
# For localhost-only dev you may allow no auth:
# Environment=ALLOW_NO_AUTH_LOCAL=1
SyslogIdentifier=legoarm
StandardOutput=journal
StandardError=journal
Restart=on-failure

[Install]
WantedBy=multi-user.target
